# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:django]
command=gunicorn project.wsgi:application --bind 0.0.0.0:8000 --workers 3 --worker-class gevent --worker-connections 1000
directory=/app/src
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/django.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:daphne]
command=daphne -b 0.0.0.0 -p 8001 project.asgi:application
directory=/app/src
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/daphne.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
